# deployment/files/templates/supervisor-site.ini
[program:{{ app_name }}]
command={{ webapps_dir }}/{{ app_name }}/venv/bin/gunicorn {{ app_name }}:app -c {{ webapps_dir }}/{{ app_name }}/gunicorn.py --env POWERSHAME_CONFIG={{ webapps_dir }}/{{ app_name }}/config.py
directory={{ webapps_dir }}/{{ app_name }}
user={{ user }}
autostart=true
autorestart=true
redirect_stderr=true
